# deploy/atlas.service
[Unit]
Description=Atlas Project Management App
After=network.target

[Service]
Type=simple
User=atlas
Group=atlas
WorkingDirectory=/opt/atlas/backend
EnvironmentFile=/opt/atlas/backend/.env
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/atlas/backend/data.db /opt/atlas/backend/uploads

[Install]
WantedBy=multi-user.target
